# com.nanoclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/ku/p/kumekay/nanoclaw/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ku/p/kumekay/nanoclaw</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ku/p/kumekay/nanoclaw/logs/nanoclaw.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ku/p/kumekay/nanoclaw/logs/nanoclaw.error.log</string>
    <key>UserName</key>
    <string>ku</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/ku</string>
    </dict>
</dict>
</plist>
